# deployments/systemd/sideplane-sidecar.service
[Unit]
Description=Sideplane Sidecar (read-only test)
Documentation=https://github.com/wucm667/sideplane
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/sideplane/sideplane-sidecar.env
ExecStart=/usr/local/bin/sideplane-sidecar --server ${SIDEPLANE_SERVER_URL} --node-id ${SIDEPLANE_NODE_ID} --state ${SIDEPLANE_STATE_PATH} --heartbeat-interval ${SIDEPLANE_HEARTBEAT_INTERVAL} --job-poll-interval ${SIDEPLANE_JOB_POLL_INTERVAL}
Restart=on-failure
RestartSec=5s
User=sideplane
Group=sideplane
StateDirectory=sideplane
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/sideplane

[Install]
WantedBy=multi-user.target
